# tpl/config/crontab
* * * * * * @ ls
* * * * * * @ php artisan list
预留一行